# mail/service/opendkim/run
#!/bin/bash
exec /usr/sbin/opendkim -f -x /etc/opendkim.conf 2>&1

