# mail/service/opensmtpd/run
#!/bin/bash
exec /usr/sbin/smtpd -dv -T lookup -f /etc/smtpd.conf
