# mail/service/rsyslogd/run
#!/bin/bash
exec /usr/sbin/rsyslogd -c3 -f /etc/rsyslog.conf >/dev/null 2>&1

